# infra/systemd/gxp-docs-monthly.timer
[Unit]
Description=Run GxP regulatory documents updater monthly

[Timer]
OnCalendar=monthly
Persistent=true
RandomizedDelaySec=1800
Unit=gxp-docs-monthly.service

[Install]
WantedBy=timers.target
